# systemd/snapshot@archive.service.d/override.conf
# Override del service genérico para 'archive'.
# Los backups mensuales pueden tardar horas en servers grandes;
# Type=oneshot sin override quedaría matado tras 90s (default de
# TimeoutStartSec). infinity deja correr hasta que termine.
[Service]
TimeoutStartSec=infinity
# Nice más alto (menos prioridad) que los snapshots diarios — un
# archive no debe interferir con operaciones normales del servidor.
Nice=15
IOSchedulingClass=idle
