# systemd/snapshot@db-archive.timer.d/override.conf
# Override del timer genérico para la instancia 'db-archive'.
# Backup diario a las 03:00 UTC con ±20min de jitter para que múltiples
# clientes no peguen al Drive a la vez.
[Timer]
OnCalendar=
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=20min
Persistent=true
