# systemd/snapshot@archive.timer.d/override.conf
# Override del timer genérico para la instancia 'archive'.
# Backup cold-storage mensual: día 1 a las 02:00 UTC, con ±1h de jitter
# para que no todos los clientes peguen al Drive a la vez.
[Timer]
OnCalendar=
OnCalendar=*-*-01 02:00:00
RandomizedDelaySec=1h
Persistent=true
