# systemd/snapshot@reconcile.timer.d/override.conf
# Override del timer genérico para la instancia 'reconcile'.
# Corre cada 30 minutos (con jitter) para drenar el buffer local a Drive
# apenas haya conectividad. Si Drive no está alcanzable, la ejecución
# es un no-op y termina rápido.
[Timer]
OnCalendar=
OnCalendar=*:0/30
RandomizedDelaySec=5min
