# systemd/snapshot@db-archive.service.d/override.conf
# Override del service genérico para 'db-archive'.
# Los dumps de DB pueden tardar mucho en bases grandes; infinity deja
# correr hasta que termine. IO scheduling idle para no competir con la
# carga normal del servidor.
[Service]
TimeoutStartSec=infinity
Nice=15
IOSchedulingClass=idle
